# deploy/moscow-lead-bot.service
[Unit]
Description=Moscow field leads Telegram bot (Arbi status buttons)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/pepperoni/repo/moscow-leads
EnvironmentFile=-/var/www/pepperoni/seo-agent.env
EnvironmentFile=-/var/www/pepperoni/moscow-leads.env
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=15
StandardOutput=append:/var/www/pepperoni/repo/moscow-leads/data/bot.log
StandardError=append:/var/www/pepperoni/repo/moscow-leads/data/bot.log

[Install]
WantedBy=multi-user.target
