# saidsultan_app/deploy/saidsultan-api.service
[Unit]
Description=Saidsultan AI Visibility API (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=__PROJECT_ROOT__
EnvironmentFile=__PROJECT_ROOT__/.env
ExecStart=__PROJECT_ROOT__/venv/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
